# iv_surface.service
[Unit]
Description=IV Surface Data Feed for SAC-DiF Trading Bot
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=yanzewu
WorkingDirectory=/Users/yanzewu/PycharmProjects/NLP_Final_Project_D
ExecStart=/usr/bin/python3 src/layers/layer1_signal_generation/iv_surface.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/Users/yanzewu/PycharmProjects/NLP_Final_Project_D

# Resource limits
LimitNOFILE=1024
MemoryLimit=256M

[Install]
WantedBy=multi-user.target