# sor_latency.service
[Unit]
Description=Smart Order Router Latency Tracker for SAC-DiF Trading Bot
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=yanzewu
WorkingDirectory=/Users/yanzewu/PycharmProjects/NLP_Final_Project_D
ExecStart=/usr/bin/python3 execution/sor_latency.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/Users/yanzewu/PycharmProjects/NLP_Final_Project_D

# Resource limits
LimitNOFILE=1024
MemoryLimit=256M

[Install]
WantedBy=multi-user.target