# systemd/corp_actions.service
[Unit]
Description=Corporate Actions Daemon
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=trading
WorkingDirectory=/opt/trading
Environment=PYTHONPATH=/opt/trading
ExecStart=/opt/anaconda3/bin/python3 services/corporate_actions_daemon.py
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=256M
CPUQuota=5%

[Install]
WantedBy=multi-user.target